# flexlab/drivers/systemd/alerts.timer
[Unit]
Description= Run alerts service every 15 min

[Timer]
OnCalendar=*:0/15
Persistent=true
Unit=alerts.service

[Install]
WantedBy=timers.target